# process-skyblue-restart.service
[Unit]
Description=Daily restart for Process SkyBlue

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart process-skyblue.service
