# process-bluesky.service
[Unit]
Description=Process BlueSky - Bluesky to X Cross-poster
# network-online.target で実際にネットワーク疎通可能になるまで待機
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=simple
# 実行ユーザーを環境に合わせて変更してください
User=your-username
# リポジトリのクローン先に合わせて変更してください
WorkingDirectory=/path/to/process_bluesky
# 残骸コンテナがあれば削除（異常終了時の名前衝突対策）
# - を付けることで失敗してもサービス起動を継続
ExecStartPre=-/usr/bin/docker rm -f process-bluesky
ExecStart=/usr/bin/docker run --rm --name process-bluesky \
    --env-file /path/to/process_bluesky/.env \
    -v /path/to/process_bluesky/data:/app/data \
    process-bluesky
ExecStop=-/usr/bin/docker stop process-bluesky
Restart=always
RestartSec=30

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=process-bluesky

[Install]
WantedBy=multi-user.target
